# data/supervisord.conf
[supervisord]
user        = root
nodaemon    = true
logfile     = /var/log/supervisor/supervisord.log
pidfile     = /var/run/supervisord.pid
childlogdir = /var/log/supervisor

[program:kubectl-proxy]
command                 = /usr/bin/kubectl proxy --address=127.0.0.1 --port=8080
user                    = prometheus
environment             = KUBECONFIG="/etc/kube/config"
autostart               = true
autorestart             = true
stdout_logfile          = /dev/stdout
stdout_logfile_maxbytes = 0
stdout_events_enabled   = true
stderr_logfile          = /dev/stderr
stderr_logfile_maxbytes = 0
stderr_events_enabled   = true

[program:transformer]
command                 = uwsgi --ini uwsgi.ini
user                    = prometheus
directory               = /home/prometheus/transform
autostart               = true
autorestart             = true
stdout_logfile          = /dev/stdout
stdout_logfile_maxbytes = 0
stdout_events_enabled   = true
stderr_logfile          = /dev/stderr
stderr_logfile_maxbytes = 0
stderr_events_enabled   = true
